# OpenVPN-Randomizer/usr/lib/systemd/system/openvpn-custom.service
[Unit]
Description=Start OpenVPN autoconnect to my vpn service 
[Service]
Type=forking
User=root
ExecStart=/opt/OpenVPN/randomize_start.sh
[Install]
WantedBy=multi-user.target
